# wlanguests/wlanguests.path
[Unit]
Description=Triggers associated service when directory contents change
Documentation=man:systemd.path

[Path]
PathModified=/var/www/wlanguests/

[Install]
WantedBy=multi-user.target
