# wlanpindisplay/wlanpindisplay.service
[Unit]
Description=Display PIN for WLAN guests
Documentation=man:systemd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wlanpindisplay
